# 2k16-Dont_be_too_Evil/supervisor_evil.conf
[program:evil]
command=python /root/evil.py

[program:dummy]
command=python /root/dummy.py
